# xdg-autostart/on-login-gui.desktop
[Desktop Entry]
Exec=~/imemine/bin/on-login-gui
Name=on-login-gui
Type=Application
Version=1.0
